# etc/systemd/system/docker.service.d/docker.conf
[Service]
Restart=always
ExecStart=
ExecStart=/usr/bin/docker daemon -g /data/docker -H fd:// $DOCKER_OPTS
