# packaging/appliance-images/openSUSE/13.2/AIO-Simulators/templates/root/AIO_scripts/unit_files/CiscoMgmt.service
[Unit]
Description=The CISCO management server
After=CISCO.service

[Service]
Type=forking
ExecStart=/etc/init.d/CiscoMgmt start
ExecStop=/etc/init.d/CiscoMgmt stop

[Install]
WantedBy=multi-user.target

