# packaging/appliance-images/openSUSE/13.2/AIO-Simulators/templates/root/AIO_scripts/unit_files/CISCO.service
[Unit]
Description=The CISCO simulator server
After=SimService.service

[Service]
Type=forking
ExecStart=/etc/init.d/CISCO start
ExecStop=/etc/init.d/CISCO stop

[Install]
WantedBy=multi-user.target
