# packaging/appliance-images/openSUSE/13.2/AIO-Simulators/templates/root/AIO_scripts/unit_files/fcnsDB.service
[Unit]
Description=The fcns.db script
After=RP.service

[Service]
Type=forking
ExecStart=/etc/init.d/fcnsDB start
ExecStop=/etc/init.d/fcnsDB stop

[Install]
WantedBy=multi-user.target



