# packaging/appliance-images/openSUSE/13.2/AIO-Simulators/templates/root/AIO_scripts/unit_files/VPLEX.service
[Unit]
Description=The VPLEX simulator
After=ECOM80.service

[Service]
Type=forking
ExecStart=/etc/init.d/VPLEX start
ExecStop=/etc/init.d/VPLEX stop

[Install]
WantedBy=multi-user.target



