# packaging/appliance-images/openSUSE/13.2/AIO-Simulators/templates/root/AIO_scripts/unit_files/VPLEX2.service
[Unit]
Description=The VPLEX2 simulator
After=ECOM80.service

[Service]
Type=forking
ExecStart=/etc/init.d/VPLEX2 start
ExecStop=/etc/init.d/VPLEX2 stop

[Install]
WantedBy=multi-user.target
