# packaging/appliance-images/openSUSE/13.2/AIO-Simulators/templates/root/AIO_scripts/unit_files/WIN.service
[Unit]
Description=The Windows simulator
After=CISCO.service

[Service]
Type=forking
ExecStart=/etc/init.d/WIN start
ExecStop=/etc/init.d/WIN stop

[Install]
WantedBy=multi-user.target
