# packaging/appliance-images/openSUSE/13.2/AIO-Simulators/templates/root/AIO_scripts/unit_files/XIO2.service
[Unit]
Description=The XIO2 simulator
After=CISCO.service

[Service]
Type=forking
ExecStart=/etc/init.d/XIO2 start
ExecStop=/etc/init.d/XIO2 stop

[Install]
WantedBy=multi-user.target


